# .config/systemd/user/restic@local.service.d/override.conf
[Service]
Environment=RESTIC_REPOSITORY=%h/.backup/restic/fedora_silverblue-39_T480
